# srcs/wordpress/srcs/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php-fpm]
command=/usr/sbin/php-fpm7 -F

[program:nginx]
command=nginx -g "daemon off;"