# srcs/nginx/srcs/supervisord.conf
[supervisor]
nodaemon=true
user=root

[program:sshd]
command=/etc/init.d/sshd start

[program:nginx]
command=nginx -g "daemon off;"